# rtorrent.ini
[program:rtorrent]
user=
command=/usr/bin/rtorrent -n -o import=/config/rtorrent.conf
dht=disable
peer_exchange=no